# deploy/roles/clickhouse/templates/clickhouse-compose.service.j2
[Unit]
Description=ClickHouse Docker Compose stack
Requires=docker.service
Wants=network-online.target tailscaled.service wg-quick@wg0.service
After=docker.service network-online.target tailscaled.service wg-quick@wg0.service

[Service]
Type=oneshot
RemainAfterExit=true
User={{ deploy_user }}
WorkingDirectory={{ deploy_dir }}
Environment="DOCKER_HOST=unix:///var/run/docker.sock"
ExecStartPre=/bin/sh -c 'until tailscale status --peers=false 2>/dev/null | grep -q "^100\\."; do sleep 1; done'
ExecStartPre=/bin/sh -c 'until ip addr show wg0 2>/dev/null | grep -q "inet 10\\.99\\."; do sleep 1; done'
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose stop
TimeoutStartSec=0
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
